# scripts/qs_host_bits2511_wow.service
# systemd --user unit for running Bluesky queueserver host for bits2511 on wow
#
# Save this file as:
#   ~/.config/systemd/user/qs_host_bits2511_wow.service
#
# How to use (no sudo required):
# 1) Reload user units:
#      systemctl --user daemon-reload
# 2) Start the service now:
#      systemctl --user start qs_host_bits2511_wow.service
# 3) Enable it to start automatically when you log in:
#      systemctl --user enable qs_host_bits2511_wow.service
# 4) Check status and logs:
#      systemctl --user status qs_host_bits2511_wow.service
#      journalctl --user -u qs_host_bits2511_wow.service -f
# 5) Stop / disable:
#      systemctl --user stop qs_host_bits2511_wow.service
#      systemctl --user disable qs_host_bits2511_wow.service
#
# Important notes and caveats:
# - This is a per-user service. You do NOT need sudo to install or run it.
# - systemd user services run while your user session is active. By default they do
#   NOT continue running after you log out. To allow the service to keep running
#   after logout ("linger"), an administrator must run:
#       sudo loginctl enable-linger <your-username>
# - Conda: systemd user environment is minimal.
#     - Source conda's activation script and activate the env (more robust).

[Unit]
Description=bits2511 QS host server (conda env: bits2511) for user jemian
After=network.target
Wants=network.target
ConditionHost=wow.xray.aps.anl.gov

[Service]
Type=simple

WorkingDirectory=/home/beams/JEMIAN/Documents/projects/prjemian/bits2511

# Additional environment variables:
Environment=CONDA_BASE=/APSshare/miniconda/x86_64
Environment=PATH=/APSshare/miniconda/x86_64/bin:/usr/local/bin:/usr/bin:/bin
Environment=COMMAND=./scripts/qs_host_bits2511.sh

ExecStart=/bin/bash -lc '"${CONDA_BASE}/bin/conda" run --live-stream -n bits2511 "${COMMAND}"'

# Other services might use these settings.
# ExecStop=
# ExecReload=

Restart=on-failure
RestartSec=15
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
